# config/includes.chroot/etc/skel/.config/autostart/d2i-welcome.desktop
[Desktop Entry]
Type=Application
Name=D2I Welcome Message
Exec=/usr/local/bin/d2i-welcome.sh
Terminal=false
