# config/includes.chroot/etc/skel/.config/autostart/thunderbird.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/thunderbird
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Thunderbird
